# debian/openisr.postinst
#!/bin/bash
set -e

# Automatically added by dh_makeshlibs
if [ "$1" = "configure" ]; then
	ldconfig
fi
# End automatically added section

update-rc.d openisr start 20 S .

if [ -x /sbin/udevcontrol ]; then
	/sbin/udevcontrol reload_rules || true;
fi

echo -e "\n You will need to run openisr-config after installing or upgrading this package."
