# debian/openisr.postinst
#!/bin/bash
set -e

# Automatically added by dh_makeshlibs
if [ "$1" = "configure" ]; then
	ldconfig
fi
# End automatically added section

# do this by hand for now
update-rc.d openisr start 20 S .

echo -e "\n You will need to run openisr-config after installing or upgrading this package."
#echo " See /usr/share/doc/openisr/README.Debian for more info."
