# debian/openisr.postinst
#!/bin/bash
set -e

# Automatically added by dh_makeshlibs
if [ "$1" = "configure" ]; then
	ldconfig
fi
# End automatically added section

# do this by hand for now
update-rc.d openisr start 20 S .

echo -e "\n You will need the openisr kernel module and a virtual machine monitor,"
echo " See /usr/share/doc/openisr/README.Debian for more info."
